# rootfs/etc/s6-overlay/s6-rc.d/api/run
#!/bin/bash
exec s6-setuidgid claude \
    s6-env HOME=/home/claude \
    node /opt/api-server/dist/index.js
